# unity7.service.in
[Unit]
Description=Unity Shell v7
Requires=unity-settings-daemon.service unity-panel-service.service bamfdaemon.service
Wants=unity-gtk-module.service
After=unity-settings-daemon.service
PartOf=graphical-session.target

[Service]
ExecStart=/usr/bin/compiz
ExecStart.Pre=@CMAKE_INSTALL_LIBDIR@/unity/unity-compiz-profile-select
Restart=on-failure
