# services/unity-panel-service.service.in
[Unit]
Description=Backing Service for the Unity Panel
After=unity7.service
PartOf=graphical-session.target

[Service]
ExecStart=${CMAKE_INSTALL_FULL_LIBDIR}/unity/unity-panel-service
Restart=on-failure
