# services/unity-panel-service.conf.in
description "Backing Service for the Unity Panel"
author "Ted Gould <ted@ubuntu.com>"

start on started unity
stop on stopped unity

emits indicators-loaded

respawn
exec ${CMAKE_INSTALL_PREFIX}/lib/unity/unity-panel-service
