# services/unity-panel-service.conf.in

start on started unity

respawn
exec ${CMAKE_INSTALL_PREFIX}/lib/unity/unity-panel-service
