# data/unity7.conf.in
description "Unity Shell v7"
author "Ted Gould <ted@ubuntu.com>, Marco Trevisan <marco@ubuntu.com>"

start on xsession SESSION=ubuntu and started unity-settings-daemon
stop on desktop-end

pre-start script
  @UNITY_LIBDIR@/upstart-prestart-check
  @UNITY_LIBDIR@/compiz-profile-selector
end script

respawn
exec compiz
