# data/unity7.service.in
[Unit]
Description=Unity Shell v7
Wants=unity-settings-daemon.service unity-gtk-module.service unity-panel-service.service bamfdaemon.service
After=unity-settings-daemon.service
PartOf=graphical-session.target

[Service]
ExecStartPre=@UNITY_LIBDIR@/systemd-prestart-check
ExecStartPre=@UNITY_LIBDIR@/compiz-profile-selector
ExecStart=/usr/bin/compiz
Restart=on-failure

[Install]
RequiredBy=ubuntu-session.target
